# config/node_server_start.conf
# /etc/init
# /etc/init.d

start on (local-filesystems)
stop on shutdown
script
  cd /.../node-app-seed/server
  exec pm2 start app.js 
end script
